# music_assistant_proxy/rootfs/etc/s6-overlay/s6-rc.d/proxy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the proxy service
# ==============================================================================

bashio::log.info "Starting Music Assistant Proxy service..."
exec /app/run.sh
