# music_assistant_proxy/rootfs/etc/s6-overlay/s6-rc.d/proxy/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the proxy service
# ==============================================================================

bashio::log.info "Starting Music Assistant Proxy..."

# Get config values
if bashio::config.has_value "server_ip"; then
    export SERVER_IP=$(bashio::config 'server_ip')
else
    bashio::log.error "Server IP is not configured"
    exit 1
fi

export SERVER_PORT=$(bashio::config 'server_port')

bashio::log.info "Configuration loaded:"
bashio::log.info "Server IP: ${SERVER_IP}"
bashio::log.info "Server Port: ${SERVER_PORT}"

# Start the proxy
exec python3 /app/proxy.py --server-ip "${SERVER_IP}" --server-port "${SERVER_PORT}"
