# music_assistant_proxy/s6-overlay/etc/services.d/proxy/finish
#!/usr/bin/with-contenv bash
declare exit_code=$1

# Cleanup
echo "Service stopped with exit code $exit_code"
