# music_assistant_proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Get config values
declare server
server=$(bashio::config 'server')

bashio::log.info "Starting Music Assistant Proxy..."
bashio::log.info "Server URL: ${server}"

# Create required directories
mkdir -p /var/log/nginx
mkdir -p /var/lib/nginx/tmp
chown nginx:nginx /var/lib/nginx/tmp

# Start Nginx
exec nginx
